# my-default-settings/target/armvirt/files/etc/uci-defaults/zz-n1-default-settings
#!/bin/sh

version=`uci -q get base_config.@status[0].version`
hostname=`uci -q get system.@system[0].hostname`
test -n "${hostname}" || hostname="OpenWrt"

	wifi_setup_radio(){
	local radio=$1
	htmode="$(uci get wireless.${radio}.htmode)"
	uci -q get wireless.${radio} >/dev/null 2>&1 && {
			uci -q del wireless.${radio}.disabled
			uci -q set wireless.${radio}.country='AU'
			uci -q set wireless.${radio}.channel='auto'
		
		if [ "$(uci -q get wireless.${radio}.band)" = "5g" ]; then
			uci -q set wireless.wifinet${radio#*o}.ssid="${SSID}_5G"
			if test -z "${htmode##HE*}"; then
				uci -q set wireless.${radio}.htmode="HE160"
			else
				uci -q set wireless.${radio}.htmode="VHT80"
			fi
		else
			uci -q set wireless.${radio}.htmode="HT40"
			uci -q set wireless.wifinet${radio#*o}.ssid="${SSID}_2.4G"
		fi
			uci -q set wireless.wifinet${radio#*o}.ieee80211r='1'
			if [ "${SSID_PASSWD}" ]; then
				uci -q set wireless.wifinet${radio#*o}.encryption='psk2'
				uci -q set wireless.wifinet${radio#*o}.key="${SSID_PASSWD}"
			else
				uci -q set wireless.wifinet${radio#*o}.encryption='none'
			fi
		}
	}
  
test $version -lt 3 && {
uci set firewall.@defaults[0].flow_offloading='1'
uci commit firewall

	[ -f /etc/config/wireless ] && {
		SSID=${hostname}
		SSID_PASSWD=""
		( while ! cat /etc/config/wireless | grep "radio0"; do
		sleep 1
		done
		for radio in radio0 radio1 radio2 radio3 wifi0 wifi1 wifi2 wifi3; do
			wifi_setup_radio ${radio}
		done
		uci commit wireless
		wifi reload
		) &
	}
}

uci -q set base_config.@status[0].version=3
uci commit base_config
